# .docker/config/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/www/supervisord.pid
logfile=/var/www/supervisord.log


[program:sshd]
command=/usr/sbin/sshd -D

[program:apache2]
command=apache2-foreground